# supervisor/cost-estimate-worker.conf
[program:cost-estimate-worker]
command=php /home/engineer/myproject/cost_estimate/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
directory=/home/engineer/myproject/cost_estimate
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/home/engineer/myproject/cost_estimate/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
